# log_gateway/rootfs/etc/services.d/log_gateway/run
#!/usr/bin/with-contenv bash
set -euo pipefail

OPTIONS_FILE="/data/options.json"

API_TOKEN="$(jq -r '.api_token // empty' "${OPTIONS_FILE}")"
Z2M_SLUG="$(jq -r '.z2m_slug // "45df7312_zigbee2mqtt"' "${OPTIONS_FILE}")"
LINES_DEFAULT="$(jq -r '.lines_default // 1000' "${OPTIONS_FILE}")"
LINES_MAX="$(jq -r '.lines_max // 1000' "${OPTIONS_FILE}")"
NO_COLORS="$(jq -r '.no_colors // true' "${OPTIONS_FILE}")"

if [[ -z "${API_TOKEN}" || "${API_TOKEN}" == "null" ]]; then
  echo "FATAL: Option 'api_token' is required." >&2
  exit 1
fi

export LOGGW_TOKEN="${API_TOKEN}"
export LOGGW_Z2M_SLUG="${Z2M_SLUG}"
export LOGGW_LINES_DEFAULT="${LINES_DEFAULT}"
export LOGGW_LINES_MAX="${LINES_MAX}"
export LOGGW_NO_COLORS="${NO_COLORS}"

echo "Starting log gateway on 0.0.0.0:8099 (snapshot only)."
exec uvicorn main:app --app-dir /app --host 0.0.0.0 --port 8099 --proxy-headers
